# code-voice/com.codevoice.stt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.codevoice.stt</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kylekillen/brainrot-radio/venv/bin/python3</string>
        <string>/Users/kylekillen/brainrot-radio/code-voice/stt_server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kylekillen/brainrot-radio</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/kylekillen/brainrot-radio/code-voice/stt.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kylekillen/brainrot-radio/code-voice/stt.launchd.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
